# sys-fs/libldm/files/libldm
#!/sbin/openrc-run

depend() {
	before localmount
	use lvm
}

command=/usr/bin/ldmtool
command_args="create all"
description="Automatically scan windows LDM volumes."

